# deploy/systemd/tutunsabri-bot.service
[Unit]
Description=TutunSabri Telegram Bot
After=network.target redis-server.service
Requires=redis-server.service
PartOf=tutunsabri.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/tutunsabri
Environment=HOME=/home/ubuntu
Environment=PATH=/home/ubuntu/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/env uv run python bot.py
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=tutunsabri.target
